# aria2-tracker-update.timer
[Unit]
Description=Aria2 tracker auto update timer

[Timer]
Persistent=true
OnUnitActiveSec=1d

[Install]
WantedBy=network.target